# res/nebula-update.service
[Unit]
Description=Auto-update for @@tun_device@@
After=nebula_@@tun_device@@.service network-online.target
Wants=network-online.target
ConditionPathExists=/etc/nebula/@@tun_device@@/node

[Service]
Type=oneshot
ExecStart=/usr/local/bin/nebula_@@tun_device@@-update.sh
User=root
StandardOutput=journal
StandardError=journal

# Exit code meanings (as defined in script):
# 0 = Update was performed successfully
# 1 = No update needed (versions match)
# 2 = Error occurred
SuccessExitStatus=0 1

[Install]
WantedBy=multi-user.target
